# chef/cookbooks/core/vdd/templates/default/varnish/upstart.conf
# start varnish on vagrant mounted
start on vagrant-mounted
exec sudo service varnish start