# chef/cookbooks/core/vdd/templates/default/nginx/upstart.conf
# start nginx on vagrant mounted
start on vagrant-mounted
exec sudo service nginx start